# scripts/phone/sorpresas.cron
# Sorpresas — calendario 2/día (valor mediodía + producto noche) + Story diaria.
# Hora de Bogotá (celular en America/Bogota). dow: 0=Dom 1=Lun ... 6=Sáb
# Publica solo si IG_PUBLISH=1 (lo controla cron_run.sh). Los reels arrancan unos
# minutos antes porque tardan ~4 min en construirse. Rutas absolutas (sin variables).

# ── Story diaria (mañana) ────────────────────────────────────────────
0 9 * * * /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh post_story.py

# ── VALOR / emocional (mediodía ~12:30) ──────────────────────────────
30 12 * * 1 /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh post_frases_carousel.py
30 12 * * 2 /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh post_guia_carousel.py
30 12 * * 3 /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh post_dedicatoria.py
30 12 * * 4 /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh post_frases_carousel.py
25 12 * * 5 /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh post_video.py --builder make_poema.py
30 12 * * 6 /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh post_guia_carousel.py
30 12 * * 0 /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh post_dedicatoria.py

# ── PRODUCTO (noche ~19:30) ──────────────────────────────────────────
25 19 * * 1 /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh post_video.py
30 19 * * 2 /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh post_carousel.py
25 19 * * 3 /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh post_video.py
30 19 * * 4 /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh post_carousel.py
25 19 * * 5 /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh post_video.py
0 19 * * 6 /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh post_image.py
25 19 * * 0 /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh post_video.py

# ── INTERACCIÓN (auto-respuestas con IA) ─────────────────────────────
*/20 * * * * /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh auto_dm.py
*/30 * * * * /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh auto_comments.py
45 8 * * * /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh auto_repost.py

# ── ANALÍTICAS / DATA (diario menciones; semanal reporte + aprendizaje) ─
30 8 * * * /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh mentions.py
50 8 * * * /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh shadowban_check.py
0 21 * * * /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh daily_stats.py
0 8 * * 1 /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh report.py
50 7 * * 1 /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh hashtag_research.py
10 8 * * 1 /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh refresh_hashtags.py
20 8 * * 1 /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh data_content.py
40 8 * * 1 /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh best_time.py

# ── SALUD del sistema (cada 3h: avisa por email si falla o cae el token) ─
0 */3 * * * /data/data/com.termux/files/home/sorpresas-worker/cron_run.sh health_check.py
